# setup/com.tallysync.manager.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tallysync.manager</string>

  <key>ProgramArguments</key>
  <array>
    <string>PLACEHOLDER_VENV_PYTHON</string>
    <string>PLACEHOLDER_MAIN_PY</string>
  </array>

  <key>WorkingDirectory</key>
  <string>PLACEHOLDER_SERVER_DIR</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <false/>

  <key>StandardOutPath</key>
  <string>PLACEHOLDER_LOG_DIR/tallysync.log</string>

  <key>StandardErrorPath</key>
  <string>PLACEHOLDER_LOG_DIR/tallysync-error.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
